# conf/supervisor_celery_beat.ini
[program:ice-celery-beat]
command=python manage.py celery beat -l INFO   ;启动命令
director=/root/cel/celery_demo                   ;指定工作目录，即manage.py所在的目录
environment=PATH="/root/.local/share/virtualenvs/cel-l64ISTw0/bin/" ;指定环境目录
stdout_logfile=/root/cel/celery_demo/logs/celery.beat.log    ;指定日志文件
stdeer_logfile=/root/cel/celery_demo/logs/celery.beat.log    ;错误日志文件
autostart=true                 ;自动启动
autorestart=true               ;自动重启
startsecs=10
stopwaitsecs=60
priority=997                   ;优先级

